# etc/cron.d/linuxmuster-cachingserver
0 */2 * * * root /usr/lib/linuxmuster-cachingserver/helper_scripts/generating-filehash-cache.sh
