# debian/postinst
#!/bin/sh
set -e

case "$1" in
    install|configure)
        exampleConfig="/var/lib/linuxmuster-cachingserver/servers.json.example"
        realConfig="/var/lib/linuxmuster-cachingserver/servers.json"

        if [ ! -f $realConfig ]; then
            cp $exampleConfig $realConfig
        fi

        mkdir -p /var/log/linuxmuster/cachingserver/

        systemctl enable isc-dhcp-server
        systemctl enable tftpd-hpa
        systemctl enable rsync

        mkdir -p /var/lib/linuxmuster/hooks/device-import.post.d
        rm -f /var/lib/linuxmuster/hooks/device-import.post.d/hook-executor-import-devices.sh
        rm -f /var/lib/linuxmuster/hooks/device-import.post.d/hack-macct-to-ad.sh

        ln -s /usr/lib/linuxmuster-cachingserver/hook-executor-import-devices.sh /var/lib/linuxmuster/hooks/device-import.post.d/hook-executor-import-devices.sh

        mkdir -p /etc/rsyncd.d

        if ! grep -Fxq "&include /etc/rsyncd.d" /etc/rsyncd.conf; then
            sed -i '1s/^/\&include \/etc\/rsyncd.d\n/' /etc/rsyncd.conf
        fi

        python3 /usr/share/linuxmuster-cachingserver/migrate-server-list.py

        systemctl restart tftpd-hpa
        systemctl restart rsync
        systemctl restart cron

        touch /var/lib/linuxmuster-cachingserver/cachingserver_rsync.secrets
        chmod 600 /var/lib/linuxmuster-cachingserver/cachingserver_rsync.secrets

        # Update WebUI permissions and restart
        python3 /usr/lib/linuxmuster-webui/etc/update_scripts/merge-permissions.py
        /usr/sbin/sophomorix-ui
        systemctl restart linuxmuster-webui

        exit 0
        ;;
    upgrade|abort-upgrade)
        echo "No work to do..."
        exit 0
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 0
        ;;
esac

#DEBHELPER#

exit 0
